# oc.service
[Unit]
Description=AMDGPU Overclock

[Service]
ExecStart=/bin/oc
WorkingDirectory=/usr/share/oc

[Install]
WantedBy=default.target
